# setup/docker/workspace/crontab/laradock
* * * * * laradock /usr/bin/php /var/www/web/artisan schedule:run >> /dev/null 2>&1